# scripts/blackbox.service
[Unit]
Description=BlackboxOS
After=docker.service
Requires=docker.service

[Service]
# Environment=DATA_DIR=/root/data
# WorkingDirectory=/root/blackbox
ExecStart=blackbox start
ExecStop=blackbox stop
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=multi-user.target