# scripts/blackbox.service
[Unit]
Description=BlackboxOS
After=docker.service
Requires=docker.service

[Service]
ExecStart=/usr/local/bin/blackbox-os start

[Install]
WantedBy=multi-user.target