# scripts/blackbox-admin.service
[Unit]
Description=BlackboxOS Admin

[Service]
ExecStart=blackbox-admin
Restart=always

[Install]
WantedBy=default.target