# scripts/blackbox-admin.service
[Unit]
Description=BlackboxOS Admin
After=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/blackbox-admin
Restart=on-failure

[Install]
WantedBy=multi-user.target